# unsorted/kindle.service/kindle-mounted.path
[Unit]
Description=Kindle is mounted

[Path]
PathExists=/home/flatcap/KINDLE

[Install]
WantedBy=paths.target

