# card-mounted/card-mounted.path
[Unit]
Description=SD card is mounted

[Path]
PathModified=/dev/mmcblk0

[Install]
WantedBy=paths.target

