# unsorted/magic-dns.service
[Unit]
Description=Magic DNS service

[Service]
ExecStart=/home/flatcap/magic-dns.pl
Restart=no

[Install]
WantedBy=default.target

