# user-timer/rich.timer
[Unit]
Description=Rich's test timer

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=default.target

