# user-timer/rich.service
[Unit]
Description=Rich's test service

[Service]
Type=oneshot
ExecStart=/home/flatcap/rich.sh

[Install]
WantedBy=default.target

